# workers/worker.conf
[program:queue-work]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:events]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=events,listeners --tries=3 --backoff=10
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:sku_events]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=sku_events,sku_listeners --tries=3 --backoff=10
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:creating-order-packing]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=creating_order_packing --tries=2 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:order-packing]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=order_packing --tries=2 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:service-price]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=service_price --tries=2 --backoff=30 --timeout=300
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:service-pack]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=service_pack --tries=2 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:m6-package-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=m6_package_event --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:shopee]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=shopee,shopee-monitor --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:kiotviet]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=kiotviet --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:lazada]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=lazada --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:tiki]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=tiki --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:tiktokshop]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=tiktokshop --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:shopbaseus]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=shopbaseus --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:sapo]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=sapo --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:topship]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=topship --tries=2 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:create-tracking-no]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=create_tracking_no --tries=2 --backoff=30
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:grant-picker]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=grant_picker --tries=2
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:cancel-tracking-no]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=cancel_tracking_no --tries=2
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:subscribe-m32-order]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan order:subscribe-m32-order
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:m32-order-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=m32_order_event --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:subscribe-fobiz-order]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan order:subscribe-fobiz-order
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:fobiz-order-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=fobiz_order_event --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:subscribe-m2-order-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan purchasing_order:subscribe-m2-order-event
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:subscribe-m32-location-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan location:subscribe-m32-location-event
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:m2-order-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=m2_order_event --tries=2
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:subscribe-m6-package-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan purchasing_package:subscribe-m6-package-event
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:m6-package-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=m6_package_event --tries=2
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:m32-location-event]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=m32_location_event --tries=2
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:calculate-warehouse-stock]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=calculate_warehouse_stock --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:document-importing]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=document_importing --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:stocks]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=stocks --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:document-freight-bill-inventory]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=document_freight_bill_inventory --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:order-calculate-amount]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=order_calculate_amount --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:calculate-balance-merchant]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=calculate_balance_merchant --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:transaction]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=transaction --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data

[program:aws]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=aws --tries=3 --backoff=30
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data


[program:order-has-inventory]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_WORKDIR)s/artisan queue:work --queue=order_calculate_has_document_inventory --tries=2
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-%(process_num)02d.log
stderr_logfile=%(ENV_APP_WORKDIR)s/storage/logs/%(program_name)s-error.log
user=www-data
